# creator-command-center.service
[Unit]
Description=Creator Command Center Web App
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/daz/.openclaw/workspace/levelup-ai/creator-app
Environment="PATH=/home/daz/.openclaw/workspace/levelup-ai/creator-app/venv/bin:/usr/bin:/bin"
Environment="YOUTUBE_API_KEY=AIzaSyA8LLamcApR0eoSr4AVv8aAwL85zZmJBoI"
Environment="PYTHONUNBUFFERED=1"
Environment="PORT=8080"
Environment="TOKEN_ENCRYPTION_KEY=daz-fixed-key-2024-secure-token-encryption-key"
ExecStart=/home/daz/.openclaw/workspace/levelup-ai/creator-app/venv/bin/python3 app.py
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
